# back/law.conf
[program:law-grade]
directory=/var/www/fastapi_app
command=/var/www/fastapi_app/venv/bin/gunicorn -k uvicorn.workers.UvicornWorker main:app -b 0.0.0.0:8000
autostart=true
autorestart=true
user=root